# tools/oxidized/tools/refresh-oxidized/refresh-oxidized.service
[Unit]
Description=Update git folder
[Service]
ExecStart=/opt/tools/refresh-oxidized.sh
