# tools/oxidized/tools/monitor/monitor.service
[Unit]
Description=Monitor syslog

[Service]
ExecStart=/opt/tools/monitor.sh

[Install]
WantedBy=multi-user.target
